# kb-kill-daemon.service
[Unit]
Description=kb-kill - disable/enable a target keyboard on a global hotkey
Documentation=https://github.com/seisinger/dotfiles
After=multi-user.target

[Service]
Type=simple
# No config path: the daemon starts config-less and is fed the active user's
# config over its control socket by each session's kb-kill-push.
ExecStart=/usr/local/bin/kb-kill-daemon run
Restart=on-failure
RestartSec=2

# --- sandboxing -------------------------------------------------------------
# Runs as root only to read /dev/input and grab keyboards; it needs nothing else.
# Config now arrives over the socket (never read from disk), and the socket is
# 0666 + SO_PEERCRED-gated (no chown), so NO capabilities are required at all -
# CAP_DAC_READ_SEARCH and CAP_CHOWN (and the @chown syscall carve-out) are gone.
RuntimeDirectory=kb-kill
RuntimeDirectoryMode=0755
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
# No network of any kind -> a compromised process cannot exfiltrate.
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
# Only input devices, nothing else under /dev.
DevicePolicy=closed
DeviceAllow=char-input rw
ProtectSystem=strict
# No home access needed any more (config is pushed, not read from ~/.config).
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0077

[Install]
WantedBy=multi-user.target
