# services/kb-kill-push.service
[Unit]
Description=kb-kill-push - feed this user's config to the kb-kill daemon
Documentation=https://github.com/seisinger/dotfiles
# Pushes config over the daemon's control socket and reconnects on its own, so
# no ordering dependency on the system daemon is needed.

[Service]
Type=simple
ExecStart=%h/.local/bin/kb-kill-push
Restart=on-failure
RestartSec=2

[Install]
# default.target (not graphical-session.target) so it also runs for pure-TTY
# logins - the daemon follows the active seat user whether graphical or on a VT.
WantedBy=default.target
